# devtarball/system/usr/share/upstart/sessions/card_parameters_setter.conf
description "Pulseaudio's droid card parameters setter"
author "Ratchanan Srirattanamet <peathot@hotmail.com>"

start on started pulseaudio
stop on session-end

respawn

exec card_parameters_setter
