# devtarball/system/etc/init/zetaw-post-boot.conf
description	"Zetaw post boot"
author		"Ratchanan Srirattanamet <peathot@hotmail.com>"

#start on started lightdm
start on desktop-session-start
#manual

task

exec setprop ctl.start zetaw-post-boot
